# osf_scraper_api/devops/templates/web_supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=/usr/local/bin/uwsgi --ini /srv/app/app.ini
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
directory=/srv/app
user=uwsgi-user

[program:worker]
command=/usr/bin/python /srv/app/osf_scraper_api/web/worker.py $OSF_QUEUE
environment=PYTHONPATH=/srv/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
directory=/srv/app
user=uwsgi-user

[program:scheduler]
command=/usr/bin/python /srv/app/osf_scraper_api/web/scheduler.py
environment=PYTHONPATH=/srv/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
directory=/srv/app
user=uwsgi-user

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0